# supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod --dbpath /data

[program:plag]
command=/plag/go-plag